# 0f-snapshot-pages@.timer
[Unit]
Description=Run 0f-snapshot-pages.service every 30 minutes

[Timer]
OnCalendar=*:0/30

[Install]
WantedBy=timers.target
